# meta-inventec/meta-cxl/recipes-phosphor/util/inventec-util/pfr-update.service
[Unit]
Description=PFR secure update

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/pfr-update
StandardOutput=journal+console
