# meta-inventec/meta-common/recipes-phosphor/gpio/phosphor-gpio-monitor/purnell-pic.service
[Unit]
Description=Purnell press power button 6 seconds assert

[Service]
Type=oneshot
ExecStart=/usr/sbin/HOTSWAP_EN.sh
